# deploy/agent.service
[Unit]
Description=SPH Network Agent Service
After=network.target backend.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/sph/restApi
Environment="PYTHONPATH=/var/www/sph/restApi"
ExecStart=/usr/bin/python3 /var/www/sph/restApi/run_agent.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

